# bootcount.service
[Unit]
Description=boot count script

[Service]
ExecStart=/data/boot-counter.sh

[Install]
WantedBy=multi-user.target
